# .linux_items/archlive/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin wktech --noclear %I 38400 linux
