# setup/linux/profile_gui/airootfs/etc/skel/.config/systemd/user/update-conky.service
[Unit]
Description=Conky config update service

[Service]
Type=oneshot
ExecStart=%h/.update_conky
