# setup/linux/include_x/airootfs/etc/skel/.config/systemd/user/update-conky.timer
[Unit]
Description=Conky config update timer

[Timer]
OnBootSec=5s
OnUnitActiveSec=30s
Unit=update-conky.service

[Install]
WantedBy=timers.target
